# deckmaster@.path
[Unit]
Description="Stream Deck Device Path"

[Path]
PathExists=/dev/%i
Unit=deckmaster@.service

[Install]
WantedBy=multi-user.target
